# roles/zabbix-server/templates/zabbix-vipd.service.j2
[Unit]
Description=Service that checks if host has a VIP attached to its interface.

[Service]
Type=oneshot
ExecStart=/bin/bash {{ keepalived_dest }}/zabbix-daemon-vip.sh
