# roles/zabbix-server/templates/zabbix-vipd.timer.j2
[Unit]
Description=Timer for zabbix-vipd.service

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s

[Install]
WantedBy=timers.target
